# scripts/com.voltstream.dartcommit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Reference copy. Install with:
    cp scripts/com.voltstream.dartcommit.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.voltstream.dartcommit.plist
  Automates the DART journal COMMIT leg daily at 16:00 ET (Mac is in America/New_York, so Hour=16 == ET).
  StartCalendarInterval => launchd runs a MISSED job on wake if the Mac was asleep at 16:00.

  TCC: the job runs via the FDA-granted DartAutoCommit.app executable (below), NOT /bin/bash, so
  macOS attributes the Full Disk Access grant to that bundle alone. A .app is used because the FDA
  picker won't accept a bare .sh. The bundle executable sources scripts/auto_commit.sh from the repo.
  Without the FDA grant, a launchd job cannot read ~/Documents ("Operation not permitted").
  Build the bundle with scripts/install_dartcommit_app.sh. See CLAUDE.md ("launchd auto-commit").
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voltstream.dartcommit</string>

    <!-- Run the .app's executable DIRECTLY so the FDA grant attaches to the DartAutoCommit.app bundle. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mikeoc/Library/Application Support/VoltStream/DartAutoCommit.app/Contents/MacOS/dart_auto_commit</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Only run at the scheduled time / on kickstart — not on every (re)load. -->
    <key>RunAtLoad</key>
    <false/>

    <!-- launchd-level capture (the script self-logs to journal/auto.log; these catch pre-exec /
         TCC failures, e.g. "Operation not permitted" if Full Disk Access is not granted).
         Both paths match the *.log gitignore rule. -->
    <key>StandardOutPath</key>
    <string>/Users/mikeoc/Documents/voltstream-ai/journal/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mikeoc/Documents/voltstream-ai/journal/launchd.err.log</string>
</dict>
</plist>
